# supervise-tor.conf
[program:tor]
command=service tor start
priority=1
redirect_stderr=true
autorestart=true